# ops/crontab
# ops/crontab — supercronic cron for the findajob scheduler container.
# Schedules mirror the systemd timers 1:1. All times are in container-local
# TZ (from env var); set per-tenant via stack .env.
#
# See docs/superpowers/specs/2026-04-17-docker-compose-design.md for the full
# mapping table.

# Ensure Python child output is line-buffered to container stdout
PYTHONUNBUFFERED=1

# ── Ingest + scoring ─────────────────────────────────────────────────────────
0    0   *  *  *   timeout ${FINDAJOB_TRIAGE_TIMEOUT:-7200} python3 /app/scripts/triage.py
*/10 *   *  *  *   timeout 900 python3 /app/scripts/watchdog.py
# Google Form ingest retired in #62 — the /ingest/ web form is the operator
# write surface now. Script kept for ad-hoc manual runs until the Google
# Form is fully decommissioned.
# */30 *   *  *  *   python3 /app/scripts/ingest_form.py

# ── Notifications ────────────────────────────────────────────────────────────
0    6   *  *  *        python3 /app/scripts/notify.py apply-reminder
15   6   *  *  *        python3 /app/scripts/notify.py daily-stats
0    7   *  *  *        python3 /app/scripts/notify.py health-check
0    8   *  *  1,3,5    python3 /app/scripts/notify.py issues-ping
# scoreboard disabled for v0.1.0 — needs `gh` CLI not present in image.
# Restore per #112 (REST API rewrite + FINDAJOB_SCOREBOARD_ENABLED gate).
# 30   8   *  *  1        python3 /app/scripts/notify.py scoreboard
0    8   *  *  0        python3 /app/scripts/notify.py feedback-review

# ── Discovery (Sunday 02:00 — competency-driven company discovery, #284) ─────
0    2   *  *  0   timeout 600 python3 /app/scripts/discover_companies.py

# ── RAG rebuild (Sunday 03:00) ───────────────────────────────────────────────
0    3   *  *  0   /usr/local/bin/aichat-ng --rag job_search_rag --rebuild-rag
